# packaging/tizenclaw-skills-secure.service
[Unit]
Description=TizenClaw Secure Skills Container
After=network.target
Before=tizenclaw.service

[Service]
Type=simple
ExecStart=/usr/libexec/tizenclaw/skills_secure_container.sh start
ExecStop=/usr/libexec/tizenclaw/skills_secure_container.sh stop
ExecReload=/usr/libexec/tizenclaw/skills_secure_container.sh restart
User=root
Group=root

[Install]
WantedBy=multi-user.target
