# packaging/tizenclaw-python.service
[Unit]
Description=TizenClaw Agent System Service (Python Port)
After=network.target network-online.target dbus.service
Requires=dbus.service

[Service]
Type=simple
User=root
Group=root
Environment=PYTHONPATH=/opt/usr/share/tizenclaw-python
ExecStart=/usr/bin/python3 /opt/usr/share/tizenclaw-python/tizenclaw_daemon.py
Restart=on-failure
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=tizenclaw-python
LimitNOFILE=65536
OOMScoreAdjust=-900

[Install]
WantedBy=multi-user.target
